# root/etc/uci-defaults/99_feiyoung
#!/bin/sh

# 配置迁移：添加新增的配置项
if ! uci -q get feiyoung.general.password_seed >/dev/null; then
    uci set feiyoung.general.password_seed=''
    uci commit feiyoung
fi

# 合并 opkg 备份的配置文件（如果存在）
if [ -f /etc/config/feiyoung-opkg ]; then
    # 临时重命名以便 uci 工具识别（uci 不支持带连字符的文件名作为配置名）
    mv /etc/config/feiyoung-opkg /etc/config/feiyoung_backup
    
    # 使用 uci 读取旧配置，完美支持多行值
    OLD_ENABLED=$(uci -q get feiyoung_backup.general.enabled)
    OLD_USERNAME=$(uci -q get feiyoung_backup.general.username)
    OLD_PWD_LIST=$(uci -q get feiyoung_backup.daily.password_list)
    
    # 恢复到新配置中
    [ -n "$OLD_ENABLED" ] && uci set feiyoung.general.enabled="$OLD_ENABLED"
    [ -n "$OLD_USERNAME" ] && uci set feiyoung.general.username="$OLD_USERNAME"
    [ -n "$OLD_PWD_LIST" ] && uci set feiyoung.daily.password_list="$OLD_PWD_LIST"
    
    uci commit feiyoung
    
    # 清理临时文件
    rm -f /etc/config/feiyoung_backup
fi

# 确保脚本有执行权限（虽然 opkg 会处理，但以防万一）
chmod +x /usr/bin/feiyoung.sh 2>/dev/null
chmod +x /usr/share/feiyoung/calc_pwd.lua 2>/dev/null

# 重载 rpcd 以应用新的 ACL 权限（忽略错误）
[ -x /etc/init.d/rpcd ] && /etc/init.d/rpcd reload 2>/dev/null

# 启用并启动服务
/etc/init.d/feiyoung enable 2>/dev/null
/etc/init.d/feiyoung restart 2>/dev/null

exit 0
